# cronjob-tmpreaper.timer
[Unit]
Description=tmpreaper cronjob timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=cronjobs.target
