# cronjob-tmpreaper.service
[Unit]
Description=tmpreaper cronjob

[Service]
ExecStart=/usr/lib/tmpreaper
