# tmpreaper.crontab
MAILTO=root

# tmpreaper to run hourly
02 * * * * root /usr/lib/tmpreaper
